# pispeed/systemd/pispeed-timer.timer
[Unit]
Description=PiSpeed Monitor - Scheduled Speedtest Timer
Documentation=https://github.com/your-username/pispeed-monitor

[Timer]
# Default: run every 60 minutes
# This is dynamically updated by the web UI when you change the interval
OnBootSec=2min
OnUnitActiveSec=60min
RandomizedDelaySec=30
Unit=pispeed-run.service

[Install]
WantedBy=timers.target
